# data/framboise-brain.service.in
[Unit]
Description=Phenix Robotik Framboise Brain software
# Wants=
# After=
DefaultDependencies=false

[Service]
Type=simple
ExecStart=@DaemonExecStart@

[Install]
WantedBy=sysinit.target
